# deploy/systemd/northstar-monitor.timer
[Unit]
Description=Run Northstar operational checks every ten minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target
