# deploy/systemd/northstar-backup.timer
[Unit]
Description=Run Northstar SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
