# deploy/systemd/northstar-data-retention.service
[Unit]
Description=Northstar privacy data retention cleanup
After=northstar.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/northstar/server
EnvironmentFile=/opt/northstar/server/.env
Environment=SUPPORT_RETENTION_DAYS=90
Environment=ANALYTICS_RETENTION_DAYS=180
Environment=OPERATIONS_RETENTION_DAYS=180
ExecStart=/usr/bin/node /opt/northstar/server/scripts/data-retention.js
