# deploy/systemd/northstar-restore-check.timer
[Unit]
Description=Verify Northstar backup restoration weekly

[Timer]
OnCalendar=Sun *-*-* 04:30:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
